# ruby/templates/etc/update-motd.d/10-uname
#!/bin/sh

# This is part of the default install, but I don't like it so commenting
#uname -snrvm