# ruby/templates/etc/update-motd.d/20-styx
#!/bin/bash

UPTIME=$(uptime -p)

cat << EOS

  (?_^
   |\ )
   |/_\      <%= v "hostname" %>

         flowing for ${UPTIME:3}

EOS